# db_resources/db_backup_resources/pg_backup.crontab
@daily /var/lib/postgresql/scripts/pg_backup_rotated.sh >> /var/lib/postgresql/backups/`date +\%Y-\%m-\%dT\%H-\%M-\%S`.log 2>&1
