# db_backup_scripts/pg_backup.crontab
@daily /var/lib/postgresql/scripts/pg_backup_rotated.sh >> /var/lib/postgresql/scripts/pg_backup.log 2>&1
